# singularity-rust.service
[Unit]
Description=Singularity Rust Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=mhugo
Group=mhugo
WorkingDirectory=/home/mhugo/code/singularity
Environment=NATS_URL=nats://127.0.0.1:4222
Environment=RUST_LOG=info
ExecStart=/home/mhugo/code/singularity/rust/target/release/singularity-rust-service
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
